# packages/openwrt/crystal
#!/bin/sh /etc/rc.common
#
# 用作 OpenWrt /etc/init.d/crystal 服务入口。

START=19
NAME=crystal
CUR_DIR=$(dirname "$(readlink -f /etc/init.d/$NAME)")

start() {
    "$CUR_DIR/ctrl" start
}

stop() {
    "$CUR_DIR/ctrl" stop
}

restart() {
    "$CUR_DIR/ctrl" restart
}
